# monitored_server/files/node_exporter.service
[Unit]
Description=Node exporter

[Service]
User=node_exporter
Group=node_exporter
ExecStart=/Node/node_exporter

[Install]
WantedBy=default.target

